# kubekey/liebe-virtual-machine/backup-etcd.service
[Unit]
Description=Backup ETCD
[Service]
Type=oneshot
ExecStart=/usr/local/bin/kube-scripts/etcd-backup.sh
